# 99-monitor.rules
ACTION=="change", SUBSYSTEM=="drm", KERNEL=="card0", RUN+="/home/user/monitor_script.sh"

